# startup/xen-guest-agent.service
[Unit]
Description=Xen guest agent
Conflicts=xe-linux-distribution.service

[Service]
ExecStart=/usr/sbin/xen-guest-agent
Restart=on-failure

[Install]
WantedBy=multi-user.target
